# mutt/emailnotifier.service
[Unit]
Description=Mutt email notifier

[Service]
ExecStart=/usr/local/bin/email_notifier
Restart=always

[Install]
WantedBy=default.target
